# ftp/assets/rootfs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon = true

[program:vsftpd]
command=/etc/init.d/vsftpd start
